# fwts-livecd-rootfs/live-build/ubuntu-cpc/includes.chroot/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
Type=oneshot
ExecStart=
ExecStart=-/bin/sh -c 'while true; do /usr/bin/fwts-frontend-text noautoshutdown; /bin/bash; done'
Restart=no
StandardInput=tty
StandardOutput=tty
